# crontab
*/10 * * * * python3 /app/epg_worker.py
0 */2 * * * python3 /app/keshet_worker.py
